# conf/etc/service/ganglia-monitor/run
#!/bin/bash
rm -f /var/run/ganglia-monitor.pid
exec /usr/sbin/gmond --pid-file=/var/run/ganglia-monitor.pid
