# conf/etc/service/gmetad/run
#!/bin/bash
rm -f /var/run/gmetad.pid
exec /usr/sbin/gmetad --pid-file=/var/run/gmetad.pid
